# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g 'daemon off;'
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true

[program:server]
command=/usr/local/bin/bun /app/server/dist/index.js
directory=/app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
environment=NODE_ENV=production,PORT=3000,INFERENCE_URL=http://localhost:8000,MODELS_PATH=/app/data/models

[program:inference]
command=/opt/inference-venv/bin/uvicorn src.main:app --host 0.0.0.0 --port 8000
directory=/app/inference
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
environment=INFERENCE_MODELS_PATH=/app/data/models,INFERENCE_CACHE_DIR=/app/data/cache/prompts,INFERENCE_DEVICE=%(ENV_INFERENCE_DEVICE)s

[program:pocketbase]
command=/usr/local/bin/pocketbase serve --http=0.0.0.0:8090 --dir=/app/db/pb_data
directory=/app/db
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true